# LifeTrac-v25/DESIGN-CONTROLLER/firmware/tractor_x8/systemd/lifetrac-logger.service
[Unit]
Description=LifeTrac v25 tractor X8 logger (MQTT → SQLite + JSONL audit)
After=network-online.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/lifetrac/firmware/tractor_x8/logger_service.py
Restart=on-failure
RestartSec=5
# IP-304: see lifetrac-camera.service.
Environment=PYTHONPATH=/opt/lifetrac/base_station

[Install]
WantedBy=multi-user.target
